# overlay/etc/supervisor/conf.d/sniproxy.conf
[program:sniproxy]
command=/usr/sbin/sniproxy -f -c /etc/sniproxy.conf 
stdout_events_enabled=true
stderr_events_enabled=true